# Modules/Protocol Benchmarks/dist-benchmarks/dist-benchmark-setup/pb-node.service
[Unit]
Description=etcd server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/fish /root/start-pb.fish

[Install]
WantedBy=multi-user.target
